# ttnmapper.conf
[program:ttnmapper]
command=/home/pi/pisupply-rak811-python-ttnmapper/ttnmapper.py
numprocs=1
redirect_stderr=false
stdout_logfile=/var/log/ttnmapper.log
autostart=true
autorestart=unexpected
